# docker/8.1/supervisord-prod.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:php]
command=php-fpm -d variables_order=EGPCS
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:requests.product_type.list]
command=/usr/local/bin/php /srv/cms/artisan amqp:work requests:product_type.list
redirect_stderr=true
stdout_logfile=/srv/cms/storage/logs/worker.log
autostart=true
autorestart=true
user=www-data
startretries=86400
startsecs=0

[program:requests.product_type.export_props]
command=/usr/local/bin/php /srv/cms/artisan amqp:work requests:product_type.export_props
redirect_stderr=true
stdout_logfile=/srv/cms/storage/logs/worker.log
autostart=true
autorestart=true
user=www-data
startretries=86400
startsecs=0

[program:events.pricelist.parsed]
command=/usr/local/bin/php /srv/cms/artisan amqp:work events:pricelist.parsed
redirect_stderr=true
stdout_logfile=/srv/cms/storage/logs/worker.log
autostart=true
autorestart=true
user=www-data
startretries=86400
startsecs=0

[program:requests.parse_rules.list]
command=/usr/local/bin/php /srv/cms/artisan amqp:work requests:parse_rules.list
redirect_stderr=true
stdout_logfile=/srv/cms/storage/logs/worker.log
autostart=true
autorestart=true
user=www-data
startretries=86400
startsecs=0
