# docker/8.1/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:php]
command=/usr/bin/php -d variables_order=EGPCS /var/www/html/artisan serve --host=0.0.0.0 --port=80
user=sail
environment=LARAVEL_SAIL="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# (!) Не запускайте vite в виртуалке, это всегда нужно запускать отдельной командой, чтобы лог был перед глазами.
# [program:vite]

[program:requests.product_type.list]
command=/usr/bin/php /var/www/html/artisan amqp:work requests:product_type.list
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker.log
autostart=true
autorestart=true
user=sail
startretries=86400

[program:requests.product_type.export_props]
command=/usr/bin/php /var/www/html/artisan amqp:work requests:product_type.export_props
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker.log
autostart=true
autorestart=true
user=sail
startretries=86400

[program:requests.parse_rules.list]
command=/usr/bin/php /var/www/html/artisan amqp:work requests:parse_rules.list
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker.log
autostart=true
autorestart=true
user=sail
startretries=86400


